# systemd/mount/local/extra.mount
[Unit]
Description=Mount Extra disk space

[Mount]
What=/dev/<DEVICE>
Where=/extra
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target
